# ansible_collections/digaxfr/system/roles/frr/templates/route-flipper.service.j2
[Unit]
Description=Default route flipper

[Service]
ExecStart=/usr/local/bin/route-flipper.sh
