# ansible_collections/digaxfr/system/roles/frr/templates/route-flipper.timer.j2
[Unit]
Description=Run Route Flipper every minute
After=frr.service

[Timer]
OnUnitActiveSec=1min
OnBootSec=1min

[Install]
WantedBy=timers.target
